# area51-20090802-200946-buildscripts.tbz/lesslinux-builder/scripts/initramfs/live-cd/etc/rc.d/0100-loopback.sh
#!/static/bin/ash

PATH=/static/bin:/static/sbin:/bin:/sbin

. /etc/rc.conf
. /etc/rc.subr/extractbootparams
. /etc/rc.subr/colors

case $1 in 
    start)
	printf "$bold===> Setting up loopback networking $normal \n"
	ifconfig lo 127.0.0.1 
    ;;
esac
