# cbsicher-area51-20100702-093049.iso/boot/isolinux/initram.img/etc/inittab
::sysinit:/etc/rc
tty1::respawn:/static/bin/show_placeholder
tty2::respawn:/static/sbin/getty -l /static/bin/login 38400 tty2
tty6::respawn:/static/bin/x_user_start
::ctrlaltdel:/static/sbin/reboot
::shutdown:/etc/rc.shutdown
