# notfall-dvd-2011-stonehenge-20110201-191153-buildscripts.tbz/lesslinux-builder/scripts/stage03/.svn/text-base/5270_lshw-gui.xml.svn-base
<llpackages>
	<package 
		name="lshw-gui"
		version="B.02.15"
		class="user"
		>
	</package>
	<scripts>
	<scriptdata location="/etc/rc.d/0512-lshw-gui.sh" mode="0755" group="0" owner="0">
		<![CDATA[#!/static/bin/ash
		
#lesslinux provides lshwgui
#lesslinux license BSD
#lesslinux description

PATH=/static/bin:/static/sbin:/bin:/sbin:/usr/bin:/usr/sbin
export PATH

. /etc/rc.defaults
. /etc/rc.subr/extractbootparams
. /etc/rc.subr/colors

case $1 in 
    start)
	printf "$bold===> Adjusting sudo for lshw-gui $normal\n"
	echo '' >> /etc/sudoers
	echo '# added by /etc/rc.d/0512-lshw-gui.sh' >> /etc/sudoers
	echo 'surfer  ALL = NOPASSWD: /usr/sbin/lshw' >> /etc/sudoers
	echo 'surfer  ALL = NOPASSWD: /usr/sbin/gtk-lshw' >> /etc/sudoers
    ;;
esac

		]]>
	</scriptdata>
	</scripts>
</llpackages>
