# lesslinux-search-and-rescue-uluru-20190720-112947.iso/boot/kernel/initram.img/etc/rc.d/9998-update.sh
#!/static/bin/ash
		
#lesslinux provides updateoverlay2
#lesslinux license BSD
#lesslinux description
# Overwrite this script at update time

PATH=/bin:/sbin:/usr/bin:/usr/sbin:/static/bin:/static/sbin
export PATH

. /etc/rc.defaults
. /etc/rc.subr/extractbootparams
. /etc/rc.subr/colors

exit 0
	
		
